# keepwsl.service
[Unit]
Description=keepwsl.service

[Service]
ExecStart=/mnt/c/Windows/System32/wsl.exe --exec sh -c "sleep infinity"
Type=simple
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target